# supervisord.conf
[supervisord]
nodaemon=true 
 
[program:mongodb] 
command=/usr/bin/mongod --port 27017 --smallfiles 
 
[program:librenews] 
command=/usr/local/bin/python -u /code/src/webserver.py
